# templates/nrpe.service.erb
[Unit]
Description=Daemon to answer plugin requests from Nagios on remote hosts
After=network.target
Requires=network.target

[Service]
Type=forking
User=<%= @c_user %>
Group=<%= @c_group %>
EnvironmentFile=-/etc/sysconfig/nrpe
ExecStart=/usr/sbin/nrpe -c <%= @c_config_dir %>/nrpe.cfg -d $NRPE_SSL_OPT


[Install]
WantedBy=multi-user.target
